# filechanges/lib/systemd/system/shairport-sync.service.d/enable.conf
[Unit]
ConditionPathExists=/var/opt/musicbox/enable_shairport
