# filechanges/etc/network/if-up.d/iptables
#!/bin/sh
iptables-restore < /etc/firewall/musicbox_iptables
