# filechanges/lib/systemd/system/dropbear.service.d/enable.conf
[Unit]
ConditionPathExists=/var/opt/musicbox/enable_ssh
